# hostd/cmd/vmsan-gateway/vmsan-gateway.service
[Unit]
Description=vmsan Gateway Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/vmsan-gateway start
Restart=on-failure
RestartSec=5
RuntimeDirectory=vmsan
RuntimeDirectoryMode=0750

[Install]
WantedBy=multi-user.target
